# refresh-geoipset.timer
[Unit]
Description=Monthly geo ipset refresh

[Timer]
#OnCalendar=monthly
OnCalendar=Tue *-*-1..7 22:00:00
RandomizedDelaySec=600
Persistent=true
 
[Install]
WantedBy=timers.target
